# data/unity8-filewatcher.conf
description "System application watcher for unity8 for Deb-based systems"
author "Michael Zanetti <michael.zanetti@canonical.com>"

start on (file FILE=/usr/share/applications/*.desktop)

exec unity8-launcher-refresh

